# deploy/local-gtm-deployer.timer
[Unit]
Description=Poll GHCR for approved Local GTM releases

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
