# deploy/local-gtm-deployer.service
[Unit]
Description=Local GTM outbound deployment poll
After=network-online.target
Wants=local-gtm-deployer.timer

[Service]
Type=oneshot
EnvironmentFile=-/etc/local-gtm/deployment.env
ExecStart=/opt/local-gtm/deploy/poller.sh
User=local-gtm-deploy
Group=local-gtm-deploy
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/local-gtm-deployer
ReadOnlyPaths=/etc/local-gtm /opt/local-gtm/deploy

[Install]
WantedBy=multi-user.target
